# etc/init.d/kbfs
#!/sbin/openrc-run

description="Keybase Filesystem service"

depend() {
    need keybase
    need keybase-redirector
}

start_pre() {
    # Set default environment variables if not already set
    : ${USER:=chronos}
    : ${UID:=1000}
    : ${HOME:=/home/${USER}}
    : ${XDG_RUNTIME_DIR:=/run/user/${UID}}

    export USER UID HOME XDG_RUNTIME_DIR

    # Backwards compatibility
    if [ -f $XDG_RUNTIME_DIR/keybase/keybase.kbfs.env ]; then
        source $XDG_RUNTIME_DIR/keybase/keybase.kbfs.env
    fi

    if [ -f $HOME/.config/keybase/keybase.autogen.env ]; then
        source $HOME/.config/keybase/keybase.autogen.env
    fi

    if [ -f $HOME/.config/keybase/keybase.env ]; then
        source $HOME/.config/keybase/keybase.env
    fi

    # Forcibly unmount /keybase in case there's anything there.
    fusermount -uz "$(keybase config get -d -b mountdir)" 2>/dev/null
}

start() {
    ebegin "Starting Keybase Filesystem service"
    start-stop-daemon --start --user ${USER} --exec /usr/bin/kbfsfuse -- -debug -log-to-file
    eend $?
}

stop() {
    ebegin "Stopping Keybase Filesystem service"
    fusermount -uz "$(keybase config get -d -b mountdir)" 2>/dev/null
    start-stop-daemon --stop --exec /usr/bin/kbfsfuse
    eend $?
}

restart() {
    svc_stop
    svc_start
}

# Specify the service should restart on failure
retry="::"
